# crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXsDEYCx installed on Tue Apr  3 11:15:36 2012)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
PATH=/home/radio:/usr/local/bin:/usr/bin:/bin:/opt/bin

#morning
00 05 * * 1-5 transcast -p swissclassic.m3u -e 2 #[Radio Suisse Classic||]
00 05 * * 0,6 transcast -p swissclassic.m3u -e 4 #[Radio Suisse Classic||]
00 07 * * 1-5 transcast -p todayfm.pls -e 2 #[The Ian Dempsey Breakfast Show|Today FM|http://www.todayfm.com/]
00 09 * * * transcast -p http://feeds.feedburner.com/radiofin -e 1 -r -t podcasts #[World Radio Choice|WFJK|http://finkeegan.com/world-radio-choice/]
00 10 * * * transcast -p http://www.rte.ie/radio1/podcast/podcast_documentaryonone.xml -e 00:40 -t podcast 
40 10 * * * transcast -p "http://www.npr.org/rss/podcast.php?id=500006" -e 00:10 -t podcast
50 10 * * * transcast -p http://feeds2.feedburner.com/tedtalks_audio/ -e 00:12 -r -t podcasts
02 11 * * * transcast -p http://www.cbc.ca/podcasting/includes/ideas.xml -e 00:58 -t podcast 

#afternoon

00 12 * * * transcast -p swissclassic.m3u -e 5
00 12 * * 0 transcast -p "http://www.npr.org/rss/podcast.php?id=510208" -e 1 -t podcast
00 12 * * 6 transcast -p http://feeds.feedburner.com/bikeshow -e 1 -t podcast
00 13 * * 6 transcast -p swissclassic.m3u -e 3
00 13 * * 0 transcast -p http://feeds.feedburner.com/wfjk -e 00:45 -t podcast
46 13 * * 0 transcast -p swissclassic.m3u -e 2:15

00 14 * * 1 transcast -p http://www.rte.ie/radio1/podcast/podcast_playback.xml 00:46 podcast
00 14 * * 2 transcast -p b006qpdd -e 00:45 -t bbc-iplayer #[Pick of the week]
00 14 * * 3 transcast -p b006slnx -e 00:30 -t bbc-iplayer #[feedback]
00 14 * * 4 transcast -p http://feeds2.feedburner.com/littleatomspodcast -e 00:30 -t podcast #[little atoms]

47 14 * * 1,2 transcast -p swissclassic.m3u -e 02:30
30 14 * * 3 transcast -p swissclassic.m3u -e 02:30

#evening
00 17 * * 1-5 transcast -p radiomalta.pls -e 02:00 -t mms

00 17 * * 3 transcast -p resonancefm.m3u -e 00:30 #[optical sound show]
00 17 * * 3 transcast -p radiomalta.pls -e 01:30 -t mms

00 17 * * 0 transcast -p http://podcasts.resonancefm.com/archives/category/shows/ost/feed -e 2 -t podcast #[OST radio show]
00 17 * * 6 transcast -p http://feeds.feedburner.com/BigScienceFM -e 1 -t podcast #[big science live show]

00 19 * * * transcast -p moomins.m3u -e 1 -r
00 20 * * 1-5 transcast -p b0100rp6 -e 3 -t bbc-iplayer #[ rathcliffe and maconie]
00 20 * * 6 transcast -p b00ds8lp -e 3 -t bbc-iplayer #[ huey show]
00 20 * * 0 transcast -p b0072ky7 -e 3 -t bbc-iplayer #[funk and soul show]

00 22 * * 0 transcast -p wbcq.pls -e 1  #[marions attic]
00 22 * * 1 transcast -p http://wfmu.org/podcast/AP.xml -e 1 -t podcast #[antique phonograph show]
00 22 * * 2 transcast -p http://www.abc.net.au/radionational/feed/2884506/podcast.xml -e 1 -t podcast #[philosophers zone]
00 22 * * 3 transcast -p http://podcasts.resonancefm.com/archives/category/shows/wavelength/feed -e 1 -t podcast #[wavelength]
00 22 * * 4 transcast -p http://www.wrn.org/listeners/world-of-radio/rss/09:00:00UTC/English/541 -e 00:30 -t podcast #[world of radio]
00 22 * * 5 transcast -p b006wr9w -e 00:30 -t bbc-iplayer #[organist entertains]
00 22 * * 6 transcast -p wbcq.pls -e 1  #[?????]
30 22 * * 1-5 transcast -p bbc_two -e 01:30 -t bbc-tv #[newsnight]

00 23 * * 0 transcast -p http://feeds.kcrw.com/kcrw/ls -e 1 -t podcast #[le show]
00 23 * * 6 transcast -p http://wfmu.org/podcast/CP.xml -e 1 -t podcast #[clay pidgeon]

25 23 * * 1-4	transcast -p otr.m3u -e 05:35 -r
00 00 * * 5,6,0 transcast -p otr.m3u -e 5 -r
00 00 * * 1 transcast -p wbcq.pls -e 1 #[al weiner workd wide]
00 01 * * 1 transcast -p wbcq-5110.pls -e 3 #[rni]
00 04 * * 1 transcast -p otr.m3u -e 1 -r
